# root/etc/s6-overlay/s6-rc.d/init-mod-universal-pip-package-install-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-universal-pip-package-install-add-package/run
